# test-sync.path
[Unit]
Description=SyncTest
[Path]
PathModified=/home/andy/test/src
MakeDirectory=true
