# 170807_path/test-path.path
[Unit]
Description=PathTest
[Path]
PathModified=/home/$USER/Desktop
