# 170807_path/test-sync.path
[Unit]
Description=SyncTest
[Path]
PathModified=/home/$USER/Desktop/src
MakeDirectory=true
