# 170807_path/test-sync.path
# This is an example file, you need to change $USER to your username
[Unit]
Description=SyncTest
[Path]
PathModified=/dev/disk/
PathExists=/media/$USER/Backup
