# 170807_path/test-timer.path
# This is an example file, you need to change $USER to your username
[Unit]
Description=Testing Paths & Timer
[Path]
PathExists=/home/$USER/Desktop/test.file
[Install]
WantedBy=multi-user.target
