# 170807_path/test-path.service
[Unit]
Description=PathTest
[Service]
ExecStart=/bin/echo "Something changed on your Desktop!"
