# 170807_path/timer-test.path
[Unit]
Description=Testing Timer Integration
[Path]
PathExists=/home/andy/test.file
